# .ddev/web-build/cronbundle.cron
* * * * * php /var/www/html/bin/console cron:run
